# vps-backup-weekly.timer
[Unit]
Description=Weekly VPS backup snapshot timer
Requires=vps-backup-weekly.service

[Timer]
OnCalendar=Sun 04:00:00 Asia/Tokyo
Persistent=true

[Install]
WantedBy=timers.target
